# config/s6/cloudflared/run
#!/usr/bin/with-contenv bash

s6-echo "Starting cloudflared"
/usr/local/bin/cloudflared proxy-dns --port 5053 --upstream https://1.1.1.1/dns-query --upstream https://1.0.0.1/dns-query